# postfixfiles/conf_d_supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/opt/nrdcs_postfix.sh

